# src/systemd/sudoku.socket
[Unit]
Description=Sudoku Solver Telnet Server

[Socket]
ListenStream=23
BindIPv6Only=both

[Install]
WantedBy=default.target
